# .config/systemd/user/set-background.timer
[Unit]
Description=~/bin/set-background.sh

[Timer]
OnCalendar=*:00:00
RemainAfterElapse=no
