# .config/systemd/user/set-background.service
[Unit]
Description=~/bin/set-background.sh

[Service]
ExecStart=%h/bin/set-background.sh
